# cleanupCart.cron
0 3 * * * cd C:\bambuEcomm/server && npx ts-node scripts/cleanupCart.ts
